# cookbooks/base/files/etc/profile.d/bash.conf.sh
alias rm='rm -i'
alias cp='cp -i'
alias mv='mv -i'
alias ls='ls -FG'
alias ll='ls -FGhl'
alias lt='ls -FGhlt'
alias la='ls -AFGhl'
