# cookbooks/base/files/etc/profile.d/less.conf.sh
export LESS=-inqMRS
export LESSCHARSET=utf-8
