# cookbooks/certbot/files/etc/cron.d/certbot
0 4 * * * root /usr/local/bin/certbot-auto renew --renew-hook "service nginx reload" --debug
